# gcp/sysroot/etc/systemd/system/cloud-init.service.d/override-GCE.conf
[Unit]
After=
After=cloud-init-local.service
After=systemd-networkd-wait-online.service
#After=google-network-daemon.service
After=google-startup-scripts.service
